# memory/secretary-signal-bot.service
# Install:
#   cp secretary-signal-bot.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now secretary-signal-bot
#
# Logs:
#   journalctl --user -u secretary-signal-bot -f

[Unit]
Description=Private Secretary Signal Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/secretary
ExecStart=%h/Projects/secretary/.venv/bin/python %h/Projects/secretary/signal_bot.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
