# memory/secretary-scheduler.service
# Install as a systemd user service so the scheduler runs on login:
#
#   mkdir -p ~/.config/systemd/user
#   cp secretary-scheduler.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now secretary-scheduler.service
#
# Check status:
#   systemctl --user status secretary-scheduler
#
# View logs:
#   journalctl --user -u secretary-scheduler -f

[Unit]
Description=Private Secretary Scheduler
After=network.target

[Service]
Type=simple
# UPDATE THESE PATHS to match your setup
WorkingDirectory=%h/Projects/secretary
ExecStart=%h/Projects/secretary/.venv/bin/python %h/Projects/secretary/scheduler.py
Restart=on-failure
RestartSec=30

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
